# .docker/crontab
#run cron at 2am
0 2 * * * root bash /backup.sh
